# com.you.mailflow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent: keeps Script 1 running so it reacts to new mail via IMAP IDLE.

  Install:
    1. Edit the paths below (python venv, project dir).
    2. cp com.you.mailflow.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.you.mailflow.plist

  Manage:
    launchctl unload ~/Library/LaunchAgents/com.you.mailflow.plist   # stop
    launchctl list | grep mailflow                                   # status
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.you.mailflow</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Use your venv's python so imap-tools is available -->
        <string>/Users/you/mailflow/.venv/bin/python</string>
        <string>/Users/you/mailflow/script1_save_attachments.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/you/mailflow</string>

    <!-- Start at login and restart if the IDLE connection drops -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/you/mailflow/mailflow.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/you/mailflow/mailflow.err.log</string>
</dict>
</plist>
